# infra/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:laravel-scheduler]
command=php /var/www/html/artisan schedule:run --verbose --no-interaction
autostart=true
autorestart=true
redirect_stderr=true
